# ansible/roles/nvim/files/nvim_cleanup.timer
[Unit]
Description=Run nvim cleanup of tags hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
