# dot_config/systemd/user/nvim-state-cleanup.timer
[Unit]
Description=Check hourly if neovim cleanup is required

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
