# dot_config/systemd/user/nvim-state-cleanup.service.tmpl
[Unit]
Description=Perform neovim cleanup

[Service]
type=oneshot
ExecStart={{ .chezmoi.homeDir }}/.local/bin/nvim-cleanup
